# example_config/systemd/etherguard-super.service
#/etc/systemd/system/etherguard-super.service
[Unit]
Description=Etherguard supernode
After=network.target

[Service]
User=root
Group=root
Type=notify
# No auxiliary socket directory is required in the current transport-agnostic runtime.
ExecStart=etherguard-go -config /etc/eggo/super.yaml -mode super
Nice=5

[Install]
WantedBy=multi-user.target
